# deploy/pos-backend.service
[Unit]
Description=Restaurant POS backend API
Documentation=https://github.com/d4not/Restaurant-POS
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=pos
Group=pos
WorkingDirectory=/opt/restaurant-pos
EnvironmentFile=/opt/restaurant-pos/.env
ExecStart=/usr/bin/node /opt/restaurant-pos/dist/server.js
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/restaurant-pos
StandardOutput=append:/var/log/restaurant-pos/backend.log
StandardError=append:/var/log/restaurant-pos/backend.err.log

# Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
